# plugins/miscellanea/playonconnect/playonconnect.service

[Unit]
Description=playonconnect service
After=syslog.target

[Service]
Type=simple
ExecStart=/bin/bash /data/plugins/miscellanea/playonconnect/playonconnect.sh
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=volumio
User=root
Group=root

[Install]
WantedBy=multi-user.target

